# deploy/personal-assistant.service
[Unit]
Description=Personal Assistant Telegram Bots (Taake & Nina)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=assistant
WorkingDirectory=/home/assistant/projekt-personal-assistent
Environment=PATH=/home/assistant/projekt-personal-assistent/venv/bin
ExecStart=/home/assistant/projekt-personal-assistent/venv/bin/python main.py
Restart=always
RestartSec=10
# Wartezeit beim Stopp (sauberes Herunterfahren)
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=personal-assistant

[Install]
WantedBy=multi-user.target
