# deploy/personal-assistant-api.service
[Unit]
Description=Personal Assistant API (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=assistant
WorkingDirectory=/home/assistant/projekt-personal-assistent
Environment=PATH=/home/assistant/projekt-personal-assistent/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/assistant/projekt-personal-assistent/venv/bin/python api/api_main.py

# Liveness-Check: prüft nach Start ob der API-Endpunkt erreichbar ist
ExecStartPost=/bin/sh -c 'for i in 1 2 3 4 5 6; do /bin/sleep 5 && /usr/bin/curl -sf http://localhost:8000/status && exit 0; done; exit 1'

# Restart with backoff
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5

# Watchdog deaktiviert: Uvicorn sendet kein sd_notify(WATCHDOG=1).
# Stattdessen prüft ExecStartPost den Health-Endpunkt nach dem Start.
# WatchdogSec=60

# Wartezeit beim Stopp (sauberes Herunterfahren)
TimeoutStopSec=30

# Ressourcen-Limits
MemoryMax=1G
MemoryHigh=768M
CPUQuota=80%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=personal-assistant-api

[Install]
WantedBy=multi-user.target
