# deploy/personal-assistant-api.service
[Unit]
Description=Personal Assistant API (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=assistant
WorkingDirectory=/home/assistant/projekt-personal-assistent
Environment=PATH=/home/assistant/projekt-personal-assistent/venv/bin
ExecStart=/home/assistant/projekt-personal-assistent/venv/bin/python api/api_main.py

# Health-Check: prüft nach Start ob der API-Endpunkt erreichbar ist
ExecStartPost=/bin/sh -c 'sleep 3 && curl -sf http://localhost:8000/health || exit 1'

# Restart with backoff
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5

# Watchdog: systemd prüft ob der Prozess noch reagiert
WatchdogSec=60

# Wartezeit beim Stopp (sauberes Herunterfahren)
TimeoutStopSec=30

# Ressourcen-Limits
MemoryMax=1G
MemoryHigh=768M
CPUQuota=80%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=personal-assistant-api

[Install]
WantedBy=multi-user.target
